# scripts/smf-dashboard.service
[Unit]
Description=SMF Dashboard (Next.js)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mikesai1
WorkingDirectory=/home/mikesai1/projects/smf-dashboard
ExecStart=/bin/bash -lc 'npm run start'
Restart=always
RestartSec=5s
Environment=NODE_ENV=production
Environment=PORT=3000
StandardOutput=journal
StandardError=journal
SyslogIdentifier=smf-dashboard

[Install]
WantedBy=multi-user.target
